# roles/fedora/templates/tomcat9.service
[Service]
ReadWritePaths=/opt /tmp /var/log/tomcat9 /var/lib/tomcat9 /usr/local/lib/fits-1.5.1 /var/lib/logs
